# slack-reminder-bot.service
[Unit]
Description=Slack Activity Reminder Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/reminder_bot_toloka
ExecStart=/opt/reminder_bot_toloka/venv/bin/python bot.py
Restart=always
RestartSec=10
StartLimitIntervalSec=0

StandardOutput=journal
StandardError=journal
SyslogIdentifier=slack-reminder-bot

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/reminder_bot_toloka/data

[Install]
WantedBy=multi-user.target
